# etc/cron.daily/cernvm-update-notification
#!/bin/sh

/usr/bin/cernvm-update-notification
